# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.log
stderr_logfile=/var/log/supervisor/sshd.log
priority=1

[program:validate]
command=python3 /workspace/validate_container.py
autostart=true
autorestart=false
startretries=0
stdout_logfile=/var/log/supervisor/validate.log
stderr_logfile=/var/log/supervisor/validate.log
priority=2

[program:pipeline-ready]
command=bash -c "echo '=== PIPELINE READY ===' > /workspace/output/READY && touch /workspace/output/DONE_MARKER 2>/dev/null; exit 0"
autostart=true
autorestart=false
startretries=0
priority=99